# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/sbin/gosu postgres /usr/lib/postgresql/15/bin/postgres -D %(ENV_PGDATA)s -c config_file=%(ENV_PGDATA)s/postgresql.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=INT

[program:backend]
command=/bin/bash /app/docker/run-backend.sh
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
